# context/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/execlineb -P
exec /usr/bin/sshd -D -e