# context/etc/s6-overlay/s6-rc.d/testcontainer-host-override/up
/bin/bash /etc/s6-overlay/s6-rc.d/testcontainer-host-override/runonce.sh